# TP14/crontab
30 02 * * * bash  ~/tp-r504/TP14/script.sh
